# .config/autostart/gnome-keyring.desktop
[Desktop Entry]
Name=GNOME Keyring
Exec=/usr/bin/gnome-keyring-daemon --start --components=secrets,ssh
Type=Application
Hidden=true
