# .config/autostart/to-do.desktop
[Desktop Entry]
Name=To-do
Exec=/home/felipec/bin/to-do
Icon=gnome-todo
Type=Application
