# .config/autostart/todo-file.desktop
[Desktop Entry]
Name=To-do list
Exec=/usr/bin/xdg-open file:///home/felipec/to-do.txt
Icon=accessories-text-editor
Type=Application
