# .config/autostart/todo-file.desktop
[Desktop Entry]
Name=To-do list
Exec=/usr/bin/xdg-open file:///home/felipec/documents/to-do.txt
Icon=gnome-todo
Type=Application
